# ansible/roles/minecraft/files/minecraft-server.socket
[Socket]
ListenFIFO=%t/minecraft-server.stdin
Service=minecraft-server.service